# docker/files/crontab
* * * * * root sms4you >> /var/log/cron.log 2>&1
* * * * * root echo "Hello world" >> /var/log/cron.log 2>&1
