# docker/files/crontab
SHELL=/bin/bash
* * * * * root cd /sms4you && python sms4you/sms4you.py >> /var/log/cron.log 2>&1
